# deployment/backup/runtime-backup.cron.example
# Run a daily runtime backup at 03:10 UTC and prune old archives.
10 3 * * * cd /opt/deeptutor && /usr/bin/python3 scripts/backup_data.py >> /var/log/deeptutor-backup.log 2>&1
20 3 * * * cd /opt/deeptutor && /usr/bin/python3 scripts/prune_backups.py --keep-daily 7 --keep-weekly 2 --keep-monthly 1 >> /var/log/deeptutor-backup.log 2>&1
